# etc/profile.d/99_icecream_path.sh
PATH="/usr/host/libexec/icecc/bin:${PATH}"
